# contrib/lumin-engine.service
[Unit]
Description=LuminAI Inference Engine
Documentation=https://github.com/kyzen7/lumin-engine
After=lumin.socket
Wants=lumin.socket

[Service]
Type=notify
ExecStart=/usr/local/bin/lumin-engine
Restart=on-failure
RestartSec=5s
StandardInput=socket
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lumin-engine
ProtectSystem=strict
ProtectHome=yes
NoNewPrivileges=true
PrivateTmp=true
RuntimeDirectory=lumin
RuntimeDirectoryMode=0755
StateDirectory=lumin
StateDirectoryMode=0700
CacheDirectory=lumin
CacheDirectoryMode=0700
ConfigurationDirectory=lumin
ConfigurationDirectoryMode=0755
Environment="LUMIN_CONFIG_DIR=%t/lumin"

# Socket activation: inherit socket from systemd
FileDescriptorStoreMax=10

[Install]
WantedBy=multi-user.target
